# src/modified/openrc/x730-pwr
#!/sbin/openrc-run

description="Monitor the Rpi GPIO for a x730 / x735 button press for reboot / shutdown."

command="/usr/bin/x730-daemon.sh"
pidfile="/run/${RC_SVCNAME}.pid"
command_background=true

depend() {

	after root
	keyword -timeout
}
